# deploy/voiceai-token.service
[Unit]
Description=Voice AI Token Server
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/voiceai
Environment=PATH=/var/www/voiceai/.venv/bin
ExecStart=/var/www/voiceai/.venv/bin/python -m uvicorn server.token_server:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
